# sources/$OEM$/$1/Software/Tasks/Upgrade.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <Triggers>
    <BootTrigger>
      <Repetition>
        <Interval>P1D</Interval>
      </Repetition>
    </BootTrigger>
    <RegistrationTrigger>
      <Repetition>
        <Interval>P1D</Interval>
      </Repetition>
    </RegistrationTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>S4U</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
    </IdleSettings>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>%ChocolateyInstall%\choco.exe</Command>
      <Arguments>upgrade all -y</Arguments>
    </Exec>
  </Actions>
</Task>